# debian/ceibal-repo-devel.cron.daily
#!/bin/sh

[ -x /usr/sbin/ceibal-update-pbuilder ] && /usr/sbin/ceibal-update-pbuilder update >> /dev/null 2>&1
