# etc/profile.d/ceibalrepo.sh
DEBEMAIL="apt@plan.ceibal.edu.uy"
DEBFULLNAME="Plan Ceibal"
export DEBEMAIL DEBFULLNAME
